# config/cbc_attach-earlyapp.service
[Unit]
Description=cbc attach
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/cbc_attach
Slice=earlyapp.slice
