# config/early-audio_resume.service
[Unit]
After=suspend.target

[Service]
ExecStart=/usr/bin/aplay /usr/share/earlyapp/beep.wav

[Install]
WantedBy=suspend.target
