# config/earlyapp_suspend.service
[Unit]
Before=suspend.target

[Service]
ExecStart=/usr/bin/systemctl stop earlyapp.service

[Install]
WantedBy=suspend.target
