# config/earlyapp.target
[Unit]
Description=earlyapp target
Wants=earlyapp.service
After=earlyapp.service

[Install]
WantedBy=basic.target
Also=earlyapp_resume.service
