# config/es2gears_resume.service
[Unit]
After=suspend.target

[Service]
ExecStart=/usr/bin/systemctl --no-block restart es2gears.service

[Install]
WantedBy=suspend.target
