# config/simple-egl_resume.service
[Unit]
After=suspend.target

[Service]
ExecStart=/usr/bin/systemctl --no-block restart simple-egl.service

[Install]
WantedBy=suspend.target
